# launchd/Library/LaunchAgents/com.servitola.gruvbox-wallpapers-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.servitola.gruvbox-wallpapers-pull</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>cd /Users/servitola/projects/gruvbox-wallpapers && git pull</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>/Users/servitola/Library/Logs/gruvbox-wallpapers-pull.err</string>
    <key>StandardOutPath</key>
    <string>/Users/servitola/Library/Logs/gruvbox-wallpapers-pull.log</string>
</dict>
</plist>
