# cron/cron_jobs/rag-improve.cron
# rag-improve: autonomous RAG corpus expansion + revisit + eval sweep
# Runs every hour; lockfile at /tmp/rag-improve.dotfiles.lock prevents overlap
# Writes to ~/projects/dotfiles/rag/{rag.eval.json, rag-eval-history.md, gaps/<collection>.md, rag-retired.md}
# Quota per run: ~15 chat calls + ~130 embed calls
0 * * * * /Users/servitola/projects/dotfiles/rag/scripts/rag-improve.py >> /tmp/rag-improve.log 2>&1
